# autokiosk/device/custom-device/device/rootfs-overlay/etc/udev/rules.d/99-hide-card.rules
KERNEL=="card2", SUBSYSTEM=="drm", ACTION=="add", RUN+="/bin/rm -f /dev/dri/card2"
